# core/conf/crontabs/crontab.ldap
# SPDX-FileCopyrightText: 2022 Zextras <https://www.zextras.com>
#
# SPDX-License-Identifier: AGPL-3.0-only

#
# crontab.ldap
#
# Monitor MDB database size
*/30 * * * * /opt/zextras/libexec/zmldapmonitordb > /dev/null 2>&1
#
