# core/conf/crontabs/crontab.logger
# SPDX-FileCopyrightText: 2022 Zextras <https://www.zextras.com>
#
# SPDX-License-Identifier: AGPL-3.0-only

#
# crontab.logger
#
# process logs
#
00,10,20,30,40,50 * * * * /opt/zextras/libexec/zmlogprocess > /tmp/logprocess.out 2>&1
#
# Daily reports
#
30 23 * * * /opt/zextras/libexec/zmdailyreport -m
